# mds/app-admin/pulse2-package-server/files/pulse2-package-server.initd
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need net mmc-agent
}

start() {
	ebegin "Starting Pulse2 Package Server"
	eval start-stop-daemon --start --quiet --exec /usr/sbin/pulse2-package-server --pidfile /var/run/pulse2-package-server.pid 
	eend $?
}

stop() {
	ebegin "Stopping Pulse2 Package Server"
	start-stop-daemon --stop --quiet --pidfile /var/run/pulse2-package-server.pid 
	eend $?
}
